# infra/systemd/sir-saathi-api.service
[Unit]
Description=SIR Saathi API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/sir-saathi
ExecStart=/opt/sir-saathi/.venv/bin/uvicorn services.api.app:create_app --factory --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5
User=sir-saathi
Group=sir-saathi
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
